# deploy/marketplace.service
[Unit]
Description=Hardin AI Marketplace Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/marketplace
Environment=PATH=/var/www/marketplace/venv/bin
ExecStart=/var/www/marketplace/venv/bin/python app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
